# humetix.service
[Unit]
Description=Humetix 입사지원 웹 애플리케이션
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/var/www/recruit
# Update the path to gunicorn as needed for your environment
ExecStart=/usr/local/bin/gunicorn --workers 3 --bind 0.0.0.0:5000 app:app
Restart=always
RestartSec=5
EnvironmentFile=/var/www/recruit/.env

[Install]
WantedBy=multi-user.target
